# vvvv.conf
[program:vvvv]
command=python3 /home/curator2/vvvv/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator2/vvvv/test.err.log
stdout_logfile=/home/curator2/vvvv/test.out.log